# dotfile/etc/profile.d/aliases.sh
# Some good standards, which are not used if the user
# creates his/her own .bashrc/.bash_profile

# --show-control-chars: help showing Korean or accented characters

# alias ls='ls -F --color=auto --show-control-chars'
# alias ll='ls -l'
alias ssh68='ssh 192.168.10.68'
alias ssh62='ssh 192.168.10.62'

######################################################################
alias dps='docker ps'
alias dims='docker images'
alias d='docker'

alias cd='cdls'

# python3
alias py='python3'
# alias p2='python2'

alias sck='shellcheck'
# systemctl
alias sysc='systemctl'
